# opt/djoe/configuration.list

etc/apache2/sites-available/proxy-djoe-docker.conf

etc/postgresql/9.3/main/postgresql.conf

etc/turnserver.conf
etc/default/coturn

lib/systemd/system/djoe.service

etc/monitorix/monitorix.conf
etc/default/sslh

opt/djoe/

opt/openfire/conf/

opt/etherpad-lite/settings.json
