# utils/configuration.list
etc/apache2/sites-available/

etc/reTurn/reTurnServer.config

etc/init.d/firewall
etc/init.d/firewall.activate
etc/init.d/firewall.deactivate

etc/init/djoe.conf
etc/init/silverpeas.conf

etc/turnserver.conf

etc/ssl/

etc/monitorix/monitorix.conf

opt/openfire/conf/openfire.xml
opt/openfire/resources/security/

opt/silverpeas-6.0-SNAPSHOT-wildfly10/export-vars.sh
opt/silverpeas-6.0-SNAPSHOT-wildfly10/bin/settings.gradle
opt/silverpeas-6.0-SNAPSHOT-wildfly10/configuration/config.properties
opt/silverpeas-6.0-SNAPSHOT-wildfly10/configuration/silverpeas/CustomerSettings.xml
opt/silverpeas-6.0-SNAPSHOT-wildfly10/properties/org/silverpeas/authentication/settings/authenticationSettings.properties

opt/etherpad-lite/settings.json
